# controler/rpizero/conf/etc/systemd/system/ifmanager.service
[Unit]
Description=Wifi manager service
After=lcd.service

[Service]
Type=simple
ExecStart=/home/pi/master-election-experiments/rpizero/wifi/manager.sh
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=piservice
User=pi
Group=pi
Restart=on-failure

[Install]
WantedBy=basic.target