# controler/rpizero/conf/etc/systemd/system/shutdown.service
[Unit]
Description=Shutdown check service
After=lcd.service

[Service]
Type=simple
ExecStart=/home/pi/master-election-experiments/rpizero/shutdown/check.sh
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=piservice
User=pi
Group=pi
Restart=on-failure

[Install]
WantedBy=basic.target